# ebooks-sync.timer
[Timer]
OnBootSec=1 sec
OnActiveSec=5 min

[Install]
WantedBy=default.target